# media-automount@.service
[Service]
Type=forking
GuessMainPID=no
ExecStart=/usr/bin/media-automount %I
