# usb-automount@.service
[Service]
Type=forking
ExecStart=/usr/bin/usb-automount %I
