# usb-automount@.service
[Service]
Type=forking
GuessMainPID=no
ExecStart=/usr/bin/usb-automount %I
